# systemd/gex-rollup-daily.timer
[Unit]
Description=gex-cron-runner: daily trading rollup at 00:05 ET
# NOTE: do NOT add Requires=gex-rollup-daily.service here — that triggers the
# service on timer activation (e.g., `systemctl start ...timer`), which is
# not what we want. The Triggers= relationship is implicit by name match.

[Timer]
# 00:05 ET every day. systemd timer interprets timezone via OnCalendar.
OnCalendar=*-*-* 00:05:00 America/New_York
# Don't fire missed runs on boot — daily rollup needs to be triggered exactly once
Persistent=false
# Random delay 0-30s to avoid the cron storm at exact midnight
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
