# systemd/gex-rollup-live-close.timer
[Unit]
Description=gex-cron-runner: trigger live-close at 16:00 ET (flip client to rollup)
# NOTE: no Requires= — that would trigger the service on timer activation.

[Timer]
OnCalendar=Mon..Fri *-*-* 16:00:00 America/New_York
Persistent=false
RandomizedDelaySec=5

[Install]
WantedBy=timers.target
